# roles/docker/templates/docker-compose@.service.j2
{{ ansible_managed | comment }}
# Boot unit for the {{ item.name }} docker compose stack. Rendered by amun-docker
# from {{ item.dir }}/resiliency.yml. Runs `docker compose up -d` at boot so the
# stack recovers regardless of Docker's per-container restart state (which is
# lost on an ungraceful reboot). Ongoing drift/exited repair is handled by the
# /etc/cron.d/{{ item.name }}-reconcile job.
{% set b = item.boot %}
{# systemd .mount unit name for the NFS path, matching `systemd-escape -p
   --suffix=mount`: strip leading/trailing slashes, escape literal '-' to
   '\x2d' (it is the path separator in the escaped name), then '/' -> '-'.
   `replace` (plain string sub) is used for the escape so re.sub never sees a
   bad '\x' escape. #}
{% set _p = (b.nfs_mount | default('')) | regex_replace('^/+', '') | regex_replace('/+$', '') %}
{% set mount_unit = (_p | replace('-', '\\x2d') | replace('/', '-')) ~ '.mount' if (b.nfs_mount is defined and b.nfs_mount) else '' %}
{% set after = ['docker.service', 'network-online.target'] + ([mount_unit] if mount_unit else []) + (b.after | default([]) | map('regex_replace', '^(.+)$', 'docker-compose@\\1.service') | list) %}
[Unit]
Description={{ item.name }} docker compose stack (amun-docker resiliency)
Requires=docker.service
After={{ after | join(' ') }}
Wants=network-online.target
{% if b.nfs_mount is defined and b.nfs_mount %}
# RequiresMountsFor implicitly adds Requires= + After= on the mount unit, so the
# stack can never start before the NFS share is mounted.
RequiresMountsFor={{ b.nfs_mount }}
{% if b.binds_to_mount | default(false) %}
# BindsTo: if the NFS mount unit itself FAILS (mount timeout — 2026-06-22),
# tear the stack down instead of leaving it bound to the empty local
# placeholder dir. `hard,_netdev` mounts don't change .mount state on a
# transient blip (I/O just blocks), so this won't flap on momentary outages.
BindsTo={{ mount_unit }}
{% endif %}
{% endif %}

[Service]
Type=oneshot
RemainAfterExit=yes
User={{ deploy_user }}
Group={{ deploy_user }}
WorkingDirectory={{ item.working_dir }}
Environment="PWD={{ item.working_dir }}"
{% if b.nfs_mount is defined and b.nfs_mount %}
ExecStartPre=/usr/bin/mountpoint -q {{ b.nfs_mount }}
# Assert it is REALLY nfs, not the empty local placeholder a failed mount leaves.
ExecStartPre=/bin/sh -c '[ "$(stat -f -c %%T {{ b.nfs_mount }})" = nfs ]'
{% for p in b.require_paths | default([]) %}
# Belt-and-suspenders: bail before docker compose sees the bind sources if a
# required subtree is missing (e.g. the NFS share mounted but not populated).
ExecStartPre=/usr/bin/test -d {{ p }}
{% endfor %}
{% endif %}
ExecStart=/usr/bin/docker compose -p {{ item.project }} up -d --remove-orphans{% if b.force_recreate | default(false) %} --force-recreate{% endif %}

ExecStop=/usr/bin/docker compose -p {{ item.project }} down

[Install]
WantedBy=multi-user.target
